# hooks/com.switchboard.monitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd: run the token-free monitor every 5 minutes. Load with:
     launchctl load ~/Library/LaunchAgents/com.switchboard.monitor.plist
     Cron alternative (crontab -e):
       */5 * * * * /ABS/PATH/TO/.venv/bin/python /ABS/PATH/TO/hooks/sb_monitor.py /ABS/PATH/TO/REPO -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.switchboard.monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABS/PATH/TO/.venv/bin/python</string>
    <string>/ABS/PATH/TO/hooks/sb_monitor.py</string>
    <string>/ABS/PATH/TO/REPO</string>
  </array>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardErrorPath</key><string>/tmp/sb_monitor.err</string>
</dict>
</plist>
